# service/celery-worker-sync.service
[Unit]
Description=WooCommerce-Odoo Celery Worker Sync
After=network.target redis.service postgresql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/home/usuario/woocommerce_odoo/microservices/admin

# Cargar variables de entorno desde .env
EnvironmentFile=/home/usuario/woocommerce_odoo/microservices/admin/.env

# Path de Python
Environment="PYTHONPATH=/home/usuario/woocommerce_odoo/microservices/admin"
Environment="TZ=America/Havana"

# Comando para ejecutar Celery Worker Sync
ExecStart=/usr/bin/python3 -m celery -A app.celery_app worker -Q sync_queue -c 4 -n sync@%h --loglevel=info

# Logs
StandardOutput=append:/var/log/woocommerce-odoo/celery-worker-sync.log
StandardError=append:/var/log/woocommerce-odoo/celery-worker-sync-error.log

# Reiniciar automáticamente si falla
Restart=always
RestartSec=5

# Detener limpiamente
KillSignal=SIGTERM
TimeoutStopSec=60

# Límites de recursos
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
