# service/celery-beat.service
[Unit]
Description=APITM Celery Beat Scheduler Service
After=network.target redis.service postgresql.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/var/www/odoo-sync

# Cargar variables de entorno desde .env (opcional)
EnvironmentFile=-/var/www/odoo-sync/.env

# Path de Python
Environment="PYTHONPATH=/var/www/odoo-sync"
Environment="TZ=America/Havana"

# Comando para ejecutar Celery Beat
ExecStart=/var/www/venv/bin/python3 -m celery -A app.celery_app.celery_app beat --loglevel=info

# Logs
StandardOutput=append:/var/log/celery-beat.log
StandardError=append:/var/log/celery-beat-error.log

# Reiniciar automáticamente si falla
Restart=always
RestartSec=5

# Detener limpiamente
KillSignal=SIGTERM
TimeoutStopSec=30

# Límites de recursos
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target