# service/celery-worker-scheduler.service
[Unit]
Description=WooCommerce-Odoo Celery Worker Scheduler
After=network.target redis.service postgresql.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/var/www/odoo-sync

# Cargar variables de entorno desde .env (opcional)
EnvironmentFile=-/var/www/odoo-sync/.env

# Path de Python
Environment="PYTHONPATH=/var/www/odoo-sync"
Environment="TZ=America/Havana"

# Comando para ejecutar Celery Worker Scheduler
ExecStart=/var/www/venv/bin/python3 -m celery -A app.celery_app worker -Q scheduler_queue -c 2 -n scheduler@%h --loglevel=info

# Logs
StandardOutput=append:/var/log/celery-worker-scheduler.log
StandardError=append:/var/log/celery-worker-scheduler-error.log

# Reiniciar automáticamente si falla
Restart=always
RestartSec=5

# Detener limpiamente
KillSignal=SIGTERM
TimeoutStopSec=60

# Límites de recursos
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
