# service/frontend.service
[Unit]
Description=WooCommerce-Odoo Frontend Build and Serve Service
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/var/www/odoo-sync-frontend

# Variables de entorno
Environment="NODE_ENV=production"
Environment="PATH=/usr/bin:/usr/local/bin:/usr/local/bin/node:/usr/local/bin/npm"

# Servir el frontend con Vite preview (debe haberse construido previamente)
# Construir manualmente antes: cd /var/www/odoo-sync-frontend && npm install && npm run build
ExecStart=/usr/bin/npm run preview -- --host 0.0.0.0 --port 3000

# Logs
StandardOutput=append:/var/log/frontend.log
StandardError=append:/var/log/frontend-error.log

# Reiniciar automáticamente si falla
Restart=on-failure
RestartSec=5

# Límites de recursos
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
